# etc/upstart/rapture.conf
description "Rapture Service"
author      "Chris Powell"

start on filesystem or runlevel [2345]
stop on shutdown

env INSTALL_DIR=/srv/rapture
exec $INSTALL_DIR/bin/rapture-app
